# deploy/bodyweight-bot.service
[Unit]
Description=BodyWeight Telegram Bot
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/bodyweight/backend
Environment="PATH=/var/www/bodyweight/backend/venv/bin"
ExecStart=/var/www/bodyweight/backend/venv/bin/python -m app.bot.main
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
